# upstart/vehicles-online.conf
# upstart/vehicles-online.conf
#
description "vehicles-online"

start on runlevel [345]
stop on starting shutdown

respawn

script
    cd /opt/vehicles-online
    echo $$ > /var/run/vehicles-online.pid

    exec su -s /bin/bash -c 'exec /opt/vehicles-online/bin/vehicles-online  -Dconfig.file=/opt/vehicles-online/conf/application.conf -Dlogger.file=/opt/vehicles-online/conf/prod-logger.xml' vehicles-online
end script

pre-start script
    rm -f /opt/vehicles-online/RUNNING_PID
end script
